# device-backup-routeros-el9/src/main/resources/usr/lib/systemd/system/backup-routeros-weekly@.timer
# systemd service for managing all backup routeros timers.

[Unit]
Description=Backup RouterOS Weekly Timer - %i

[Timer]
OnCalendar=Sun *-*-* 01:00:00
Persistent=true
Unit=backup-routeros@%i.service

[Install]
WantedBy=timers.target

