# device-backup-rsync-el9/src/main/resources/usr/lib/systemd/system/backup-rsync-autodiscovery.service
# systemd service for managing all backup rsync timers.

[Unit]
Description=Backup Rsync Timers Autodiscovery

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/libexec/device-timer/backup-rsync start
ExecReload=/bin/bash /usr/libexec/device-timer/backup-rsync reload
RemainAfterExit=on

[Install]
WantedBy=multi-user.target

